# deb/photoboth_1_all/DEBIAN/postinst
mkdir /home/pi/data
chown -R pi:pi /home/pi/data
systemctl enable photoboth.service
