# config/picam.service
[Unit]
Description=picam

[Service]
ExecStart=raspivid -op 100 -t 0

[Install]
WantedBy=multi-user.target
